# imagens/phpbb-caddy/supervisord.conf
[supervisord]
nodaemon=true

[program:caddy]
command=caddy run --config /etc/caddy/Caddyfile --adapter caddyfile

[program:php-fpm81]
command=php-fpm81 -F